# bin/init.d/ubuntu/wcrs-address
#!/bin/sh
#
# /etc/init.d/wcrs-address -- startup script for wcrs-address
#
### BEGIN INIT INFO
# Provides:          wcrs-address
# Required-Start:    $network $remote_fs $named
# Required-Stop:     $network $remote_fs $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts wcrs-address
# Description:       Starts wcrs-address using start-stop-daemon
### END INIT INFO

case $1 in
    start)
        su - wcrs-address /caci/deploys/wcrs-address/live/bin/start.sh
    ;;
    stop)
        su - wcrs-address /caci/deploys/wcrs-address/live/bin/stop.sh
    ;;
    restart)
        su - wcrs-address /caci/deploys/wcrs-address/live/bin/stop.sh
        su - wcrs-address /caci/deploys/wcrs-address/live/bin/start.sh
    ;;
esac

exit 0

